# files/rc.local
#!/bin/bash
# =============================================================================
# Haxinator rc.local
# =============================================================================
# USB0 network configuration has been moved to haxinator-usb0-network.service
# for proper systemd dependency management.
#
# This file is kept for any future boot-time customizations.
# =============================================================================

exit 0
